# packaging/postinst
#!/bin/sh
# -------------------------------------------------------------------------------
# Post-install script for vault-cert-manager
# -------------------------------------------------------------------------------
set -e

if [ -d /run/systemd/system ]; then
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

echo "vault-cert-manager installed successfully."
echo "Configure /etc/vault-cert-manager/config.yaml and start with:"
echo "  systemctl enable --now vault-cert-manager"
